# user/window-manager.target
[Unit]
Description=Window Manager
Wants=window-manager.service
After=window-manager.service
AllowIsolate=true

[Install]
WantedBy=default.target
